# service-generator/template.service
[Unit]
Description=desc

[Service]
ExecStart=executable
WorkingDirectory=wrkdir
Restart=always

[Install]
WantedBy=multi-user.target
